# server/charging-prices/charging-prices-fetch.service
# Oneshot fetcher, triggered by charging-prices-fetch.timer.
[Unit]
Description=Wegwiesel charging-prices fetch (Mobilithek AFIR ad-hoc prices)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=wegwiesel
WorkingDirectory=/opt/wegwiesel/charging-prices
Environment=MOBILITHEK_CERT=/opt/wegwiesel/charging-prices/mobilithek.crt
Environment=MOBILITHEK_KEY=/opt/wegwiesel/charging-prices/mobilithek.key
Environment=SUBSCRIPTIONS=/opt/wegwiesel/charging-prices/subscriptions.txt
Environment=OUT_FILE=/opt/wegwiesel/charging-prices/data/charging_prices.json
ExecStart=/usr/bin/python3 /opt/wegwiesel/charging-prices/fetch_prices.py
