# server/charging-prices/charging-prices-fetch.timer
# Runs the fetcher daily (the AFIR ad-hoc prices change at most ~daily).
# Install both .service and .timer to /etc/systemd/system/, then
#   systemctl daemon-reload && systemctl enable --now charging-prices-fetch.timer
[Unit]
Description=Daily Wegwiesel charging-prices fetch

[Timer]
OnCalendar=*-*-* 04:30:00
Persistent=true
RandomizedDelaySec=900

[Install]
WantedBy=timers.target
